# supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/dev/null
logfile_maxbytes=0

; ── Flask / Python bot ────────────────────────────────────
[program:flask]
command=sh -c "gunicorn --bind 0.0.0.0:${PORT:-5000} --workers 1 --timeout 120 --access-logfile - botapp:app"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── WhatsApp bridge (Node.js) ─────────────────────────────
[program:bridge]
command=node /app/whatsapp-bridge/index.js
directory=/app/whatsapp-bridge
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
